# ansible/roles/rke2-prepare/templates/rke2-longhorn-agent.service.j2
# rke2-longhorn-agent.service.j2
[Unit]
Description=RKE2 Longhorn Agent
After=network.target

[Service]
ExecStart=/usr/local/bin/rke2 agent
KillMode=process
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target